# roles/wiki/files/supervisor.conf
[program:wiki]
command=/usr/bin/tiddlywiki /srv/wiki --server
user=wiki
autostart=true
autorestart=true
stderr_logfile=/var/log/wiki.err.log
stdout_logfile=/var/log/wiki.out.log

